# ros_noetic_catkin_ws/devel/etc/catkin/profile.d/05.catkin_make_isolated.bash
/home/jstefan1/Documents/CodeWorkspace/UGV-Robot/ros_noetic_catkin_ws/devel/.private/catkin/etc/catkin/profile.d/05.catkin_make_isolated.bash